# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ─── Loki ────────────────────────────────────────────────────────────────────
[program:loki]
command=/usr/local/bin/loki -config.file=/etc/loki/local-config.yaml
autostart=true
autorestart=true
priority=10
user=telemetry
stdout_logfile=/var/log/supervisor/loki.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/loki.err
stderr_logfile_maxbytes=5MB

# ─── Prometheus ──────────────────────────────────────────────────────────────
[program:prometheus]
command=/usr/local/bin/prometheus
    --config.file=/etc/prometheus/prometheus.yml
    --storage.tsdb.path=/data/prometheus
    --storage.tsdb.retention.time=30d
    --web.enable-remote-write-receiver
autostart=true
autorestart=true
priority=10
user=telemetry
stdout_logfile=/var/log/supervisor/prometheus.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/prometheus.err
stderr_logfile_maxbytes=5MB

# ─── Tempo ───────────────────────────────────────────────────────────────────
[program:tempo]
command=/usr/local/bin/tempo -config.file=/etc/tempo.yaml
autostart=true
autorestart=true
priority=10
user=telemetry
stdout_logfile=/var/log/supervisor/tempo.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/tempo.err
stderr_logfile_maxbytes=5MB

# ─── Grafana ─────────────────────────────────────────────────────────────────
[program:grafana]
command=/usr/share/grafana/bin/grafana server
    --homepath=/usr/share/grafana
    --config=/etc/grafana/grafana.ini
    cfg:default.paths.data=/data/grafana
    cfg:default.paths.logs=/var/log/supervisor
    cfg:default.paths.plugins=/data/grafana/plugins
    cfg:default.paths.provisioning=/etc/grafana/provisioning
autostart=true
autorestart=true
priority=20
startsecs=5
user=telemetry
stdout_logfile=/var/log/supervisor/grafana.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/grafana.err
stderr_logfile_maxbytes=5MB

# ─── Alloy ───────────────────────────────────────────────────────────────────
[program:alloy]
command=/usr/local/bin/alloy run /etc/alloy/config.alloy
    --storage.path=/data/alloy
autostart=true
autorestart=true
priority=20
user=telemetry
stdout_logfile=/var/log/supervisor/alloy.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/alloy.err
stderr_logfile_maxbytes=5MB

# ─── Nginx ───────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/nginx.err
stderr_logfile_maxbytes=5MB